# scripts/systemd/finance-bot.service
[Unit]
Description=Personal Finance OS Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=__REPO_ROOT__/.env
ExecStart=__UV_BIN__ run python -m finance_bot.main
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
