# scripts/launchd/com.personalfinance.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.personalfinance.api</string>
  <key>WorkingDirectory</key>
  <string>/path/to/Personal-Finance-OS</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/Personal-Finance-OS/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>finance_api.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/finance-api.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/finance-api.err</string>
</dict>
</plist>
